# agent_bridge/com.agentbridge.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agentbridge.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/alep/miniconda3/bin/python3</string>
        <string>-m</string>
        <string>agent_bridge.bridge_server</string>
        <string>--port</string>
        <string>8766</string>
        <string>--host</string>
        <string>127.0.0.1</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/alep/Downloads/windsurf-smoke</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/alep/Downloads/windsurf-smoke/agent_bridge/data/logs/bridge_server.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/alep/Downloads/windsurf-smoke/agent_bridge/data/logs/bridge_server.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/alep/miniconda3/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONPATH</key>
        <string>/Users/alep/Downloads/windsurf-smoke</string>
    </dict>
</dict>
</plist>
